# ansible/lets-encrypt-cleanup.service
[Unit]
Description=Let's Encrypt cleaup old certificates

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/lets-encrypt-cleanup.sh