# ansible/lets-encrypt-create.timer
[Unit]
Description=Let's Encrypt create new certificates every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target